# docker/app.crontab
# min   hour    day     month   weekday command
*   *   *   *   *   /srv/bin/console app:balance-statistic
